# minebd/systemd/nbd@nbd0.service.d/mount.conf
[Service]
ExecStartPost=/usr/lib/minebd/mount-nbd0.sh
ExecStop=/usr/bin/umount /dev/nbd0

